# gwless.service
[Unit]
Description=Gwless — Network Client Dashboard
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/gwless
ExecStart=/usr/bin/python3 -m uvicorn backend.main:app --host 0.0.0.0 --port 8080 --log-level info
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
